# com.akshay.assistant.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.akshay.assistant</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/caffeinate</string>
        <string>-i</string>
        <!-- Replace with your Python path: run `which python3` to find it -->
        <string>/Library/Frameworks/Python.framework/Versions/3.11/bin/python3</string>
        <!-- Replace YOUR_USERNAME with your macOS username -->
        <string>/Users/YOUR_USERNAME/MyAssist/personal-assistant/bot_server.py</string>
    </array>
    <key>WorkingDirectory</key>
    <!-- Replace YOUR_USERNAME with your macOS username -->
    <string>/Users/YOUR_USERNAME/MyAssist/personal-assistant</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/YOUR_USERNAME/MyAssist/personal-assistant/assistant.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/YOUR_USERNAME/MyAssist/personal-assistant/assistant.log</string>
</dict>
</plist>
